# inspircd/rootfs/etc/services.d/inspircd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the InspIRCd service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# shellcheck shell=bash

bashio::log.info "Starting InspIRCd..."

# Run inspircd
/inspircd/inspircd start

bashio::log.info "Viewing InspIRCd..."

# Get status
/inspircd/inspircd status

bashio::log.info "Stopping InspIRCd..."

# Just for testing
sleep infinity